# cli.ini
[program:cli]
command=/app/Cli start accept --token +bEIsMXbTOxfjxQWwrx96HSWwLWQJ1cbIC2Rzmepcis= --device-name DEVICE_NAME
directory=/app
autostart=true
autorestart=true
startretries=3
stderr_logfile=/var/log/cli.err.log
stdout_logfile=/var/log/cli.out.log
user=root
environment=DOTNET_SYSTEM_GLOBALIZATION_INVARIANT="0"
